# packages/blaueis-gateway/systemd/blaueis-gateway.target
[Unit]
Description=All Blaueis Gateway Instances

[Install]
WantedBy=multi-user.target
